# deploy/skrapp-worker.service
[Unit]
Description=Skrapp Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/skrapp
EnvironmentFile=/etc/skrapp/skrapp.env
ExecStart=/opt/skrapp/.venv/bin/python run_worker.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
